# backup_manager/backup.cron
0 10 * * * backupreader /backup_mysql.sh 2>> /home/backupreader/backups/backup.log > /dev/null
