# webserver/html/ropewiki/extensions/SphinxSearch/sphinx.cron
*/30 * * * * /usr/bin/indexer --config /etc/sphinxsearch/sphinx.conf --all --rotate
