# webserver/html/ropewiki/extensions/SphinxSearch/sphinx.cron
0 10 * * * /usr/bin/indexer --rotate --config /etc/sphinxsearch/sphinx.conf --all
