# database/backup.cron
0 10 * * * /bin/sh /home/root/backup_mysql.sh > /dev/null 2>&1
