# database/backup.cron
0 10 * * * bash /home/root/backup_mysql.sh >> /home/backupreader/backups/backup.log 2>&1
