# backup_manager/client/sync_backups.cron
15 10 * * * bash /root/sync_backups.sh >> /root/backup.log 2>&1
